# ddcd.service
[Install]
WantedBy=multi-user.target

[Unit]
Description=DDC daemon
Requires=ddcd.socket

[Service]
User=ddcd
Group=nobody
ExecStart=/usr/bin/ddcd
