# ddcd.socket
[Unit]
Description=DDC daemon socket

[Socket]
ListenStream=/run/ddcd/socket
